# meta-openbmc-mods/meta-common/recipes-core/kernel-panic-check/kernel-panic-check/kernel-panic-check.service
[Unit]
Description=Check for kernel panic

[Service]
Type=oneshot
ExecStart=/usr/bin/kernel-panic-check.sh

[Install]
WantedBy=multi-user.target
