# meta-openbmc-mods/meta-common/recipes-core/cpu-log-util/files/com.intel.CpuDebugLog.service
[Unit]
Description=Intel BMC CPU Debug Log

[Service]
Restart=always
ExecStart={bindir}/cpu_log
Type=simple

[Install]
WantedBy=multi-user.target
