# meta-openbmc-mods/meta-common/recipes-phosphor/selftest/intel-self-test/xyz.openbmc_project.selftest.service
[Unit]
Description= BMC Self-Test

[Service]
Restart=always
ExecStart=/usr/bin/env selftest
SyslogIdentifier=selftest

[Install]
WantedBy=multi-user.target
