# meta-openbmc-mods/meta-common/recipes-core/dhcp-check/dhcp-check/dhcp-check.service
[Unit]
Description=Check for DHCP address
After=network.target

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/bin/dhcp-check.sh'

[Install]
WantedBy=multi-user.target
