# meta-openbmc-mods/meta-common/recipes-phosphor/watchdog/system-watchdog/watchdog-reset@.service
[Unit]
Description=Reset BMC Using Hardware Watchdog

[Service]
ExecStart=/usr/bin/watchdog-reset.sh %i
Type=oneshot

