# meta-openbmc-mods/meta-common/recipes-phosphor/watchdog/system-watchdog/watchdog-clear-failures.timer
[Unit]
Description=Starts the clear watchdog serivce after 30 minutes

[Timer]
OnBootSec=30min

[Install]
WantedBy=timers.target
