# meta-openbmc-mods/meta-common/recipes-intel/configure-usb-c/files/configure-usb-c.service
[Unit]
Description=Configure USB Type C controller

[Service]
ExecStart=/usr/bin/configure-usb-c.sh
Type=oneshot

[Install]
WantedBy=basic.target
