# meta-openbmc-mods/meta-common/recipes-intel/chassis/intel-chassis-control/obmc-host-shutdown@.target
[Unit]
Description=Power%i Host Off
Wants=multi-user.target
RefuseManualStop=yes
OnFailureJobMode=flush
