# meta-openbmc-mods/meta-common/recipes-intel/chassis/intel-chassis-control/xyz.openbmc_project.Chassis.Control.Chassis@.service
[Unit]
Description=Intel Chassis%i Control

[Service]
Restart=always
ExecStart=/usr/bin/env chassis-control
SyslogIdentifier=chassis-control
Type=dbus
BusName={BUSNAME}
Nice=19
TimeoutStartSec=180s

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}

