# meta-openbmc-mods/meta-common/recipes-phosphor/fru/default-fru/SetBaseboardFru.service
[Unit]
Description=Check for FRU presence

[Service]
ExecStart=/usr/bin/checkFru.sh
Type=oneshot

[Install]
WantedBy=basic.target
