# meta-openbmc-mods/meta-common/recipes-phosphor/multi-node-nl/multi-node-nl/multi-node-nl.target
[Unit]
Description=Target for non-legacy node in multi-node system
Documentation=man:systemd.special(7)

