# meta-openbmc-mods/meta-common/recipes-phosphor/multi-node-nl/multi-node-nl/nonLegacyNode.service
[Unit]
Description=Non Legacy node

[Service]
ExecStart=/usr/bin/nonLegacyNode.sh
Type=exec

[Install]
WantedBy=multi-node-nl.target
