# meta-openbmc-mods/meta-ast2500/recipes-phosphor/fru/default-fru/SetBaseboardFru.service
[Unit]
Description=Reads GPIO and Loads the FRU

[Service]
ExecStart=/usr/bin/checkFru.sh
Type=oneshot

[Install]
WantedBy=basic.target
