# meta-openbmc-mods/meta-common/recipes-core/security-registers-check/security-registers-check/security-registers-check.service
[Unit]
Description=Check for security registers

[Service]
Type=oneshot
ExecStart=/usr/bin/security-registers-check.sh
Nice=5

[Install]
WantedBy=multi-user.target
